# test/test-app/config/supervisord.conf
[supervisord]
nodaemon = true

[program:grafana]
command = /opt/grafana/bin/grafana-server -homepath /opt/grafana/ web
autorestart = true

[program:influxdb]
command = /usr/bin/influxd
autorestart = true
